# scripts/systemd/devbox-backup.timer.template
[Unit]
Description=Daily DevBox backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1800
Unit=devbox-backup.service

[Install]
WantedBy=timers.target
